# tools/play.service
[Unit]
Description=My Play Module
After=network.target

[Service]
ExecStart=/caminho/completo/para/python /caminho/completo/para/play.py

[Install]
WantedBy=multi-user.target
