# tools/start.service
[Unit]
Description=Play Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/suehirokawashima/Documentos/trojan/play.py
WorkingDirectory=/home/suehirokawashima/Documentos/trojan/
Restart=always
User=suehirokawashima

[Install]
WantedBy=default.target
