# tools/start.service
[Unit]
Description=Games
After=network.target

[Service]
User=%i
WorkingDirectory=/chrome
ExecStart=/chrome/venv/chrome/bin/python "$(dirname "$(realpath "$0")")/modules/play.py"
Restart=always

[Install]
WantedBy=multi-user.target
